# systemd/xuance-sub.service
[Unit]
Description=XuanCe Subscription Server
After=network-online.target xray.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
ExecStart=/usr/bin/python3 /root/xuance_sub.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
